# ai/deploy/jltamp-ai.service
[Unit]
Description=JLTamp AI — semantic playlist engine
Documentation=https://github.com/JDJelectronics/JLTamp-server/blob/main/ai/README.md
After=network-online.target jltamp-embed.service
Wants=network-online.target
Requires=jltamp-embed.service

[Service]
Type=simple
User=USER
Group=USER
WorkingDirectory=/home/USER/jltamp-ai
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -m app.main
Restart=on-failure
RestartSec=10

# The service reads .env (mode 600) and writes only its own data dir. Nothing
# else on the filesystem needs to be reachable, so nothing else is.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/USER/jltamp-ai/data /home/USER/jltamp-ai/logs
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
